# src/cron.d/imscp_roundcube
# CRON(8) configuration file for RoundCube Webmail - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
#  Roundcube database cleaning: finally removes all records that are
#  marked as deleted.

0 5 * * * {USER} test -x {GUI_ROOT_DIR}/vendor/imscp/roundcube/roundcubemail/bin/cleandb.sh && {GUI_ROOT_DIR}/vendor/imscp/roundcube/roundcubemail/bin/cleandb.sh >/dev/null
