# iMSCP/cron.d/imscp_roundcube
# CRON(8) configuration file for RoundCube Webmail - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
#  Roundcube database cleaning: finally removes all records that are
#  marked as deleted.

0 5 * * * {PANEL_USER} test -x {GUI_PUBLIC_DIR}/tools/webmail/bin/cleandb.sh && {GUI_PUBLIC_DIR}/tools/webmail/bin/cleandb.sh >/dev/null
